# supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
priority=10

[program:vite]
command=npm run dev
directory=/var/www/html
autostart=true
autorestart=true
priority=20
